# packages.list
# --- MountNAS code (locally built apk) ---
mountnas-tools

# --- core shell / base utilities ---
# busybox-extras: the httpd applet serving the 'nas web' status dashboard
busybox-extras
bash
zsh
coreutils
findutils
util-linux
less
file
tree
pv
mc
tmux
nano
lsof
psmisc
jq
yq
tzdata

# --- device manager (eudev, replaces mdev) ---
eudev
udev-init-scripts
udev-init-scripts-openrc

# --- networking / transfer ---
curl
wget
rsync
rclone
# restic: encrypted, deduplicated, VERSIONED backups (rclone syncs, restic
# backs up) to local paths, SFTP, S3, or any rclone remote. Static binary.
restic
# borgbackup: same job as restic (dedup/encrypted/versioned archives) with a
# different repo format + ecosystem (borgmatic, Vorta) — both ship so users
# keep whichever workflow they already have. python3 is on the image anyway.
borgbackup
openssh
openssh-client
openssh-sftp-server
# mosh: roaming, low-latency remote shell over UDP (60000-61000). No daemon —
# mosh-server is spawned per-session over SSH, so nothing to enable.
mosh
# ttyd: browser-based terminal (the 'nas ttyd' service, OFF by default) —
# serves a real /bin/login prompt over HTTP on the trusted LAN, port 22222
ttyd

# --- mesh VPN (services OFF by default) ---
tailscale
# zerotier-one is built locally (see the /zerotier-one APKBUILD): it repackages
# the upstream prebuilt apk under our signing key, since it is not in Alpine.
zerotier-one
# plain kernel WireGuard (wg + wg-quick) for point-to-point/site-to-site
# links — the module already ships in linux-lts; no service, bring up per use
wireguard-tools

# --- host firewall (ufw — shipped DISABLED, the image is fully open) ---
# ufw drives iptables/ip6tables (one merged 'iptables' pkg in v3.24, already
# present as a docker-engine dep; python3 likewise already on the image, so
# these two cost ~0.3 MB). Rules + the enable flag live under /etc/ufw and
# persist via the normal lbu overlay. The ufw service sits in the boot
# runlevel but is a quiet no-op until 'ufw enable' (see conf.d/ufw in
# genapkovl) — enabling is one command + nas commit, no runlevel work.
ufw
ufw-openrc

# --- name resolution / discovery ---
avahi
# avahi-tools: avahi-resolve/avahi-browse — lets users (and the QEMU suite)
# actually verify <hostname>.local resolution, not just that the daemon runs
avahi-tools
dbus

# --- file sharing servers ---
samba
samba-client
samba-common-tools
nfs-utils

# --- filesystems ---
e2fsprogs
xfsprogs
btrfs-progs
f2fs-tools
exfatprogs
dosfstools
# NTFS: the in-kernel ntfs3 driver (loaded via /etc/modules) handles read/write and
# is much faster than the old FUSE ntfs-3g, so the FUSE driver is intentionally NOT
# installed. Mount NTFS disks with fstype 'ntfs3'. We keep ntfs-3g-progs for the
# userspace tools the kernel driver lacks (mkfs.ntfs, ntfsfix, ntfsresize, ...).
ntfs-3g-progs
hfsprogs
udftools
fuse
fuse3
fuse-overlayfs

# --- parity / volume mgmt (snapraid + mergerfs are built locally; see the
#     /snapraid and /mergerfs APKBUILDs — compiled from source / repackaged
#     from upstream's static binary, since neither is in latest-stable) ---
snapraid
mergerfs
mdadm
lvm2

# --- partitioning ---
parted
gptfdisk
sgdisk
cfdisk
sfdisk
# cmkfs: guided TUI for mkfs (built locally from ethanpil/cmkfs — always the
# latest upstream release, repackaged under our signing key; see /cmkfs)
cmkfs

# --- disk health / recovery / benchmarking ---
smartmontools
nvme-cli
hdparm
lsscsi
sg3_utils
ddrescue
# testdisk: partition-table recovery + PhotoRec undelete — recovery tools must
# already be on the box when the disaster happens (same logic as ddrescue)
testdisk
# f3: detect counterfeit/failing flash media (f3probe, f3read/f3write) — this
# project lives on USB sticks, so it should be able to test them
f3
fio

# --- file integrity / dedup ---
# xxhash: very fast checksums (xxhsum) for verifying large media copies
xxhash
# fdupes: find duplicate files after consolidating old drives onto the NAS
fdupes

# --- compression (busybox only DEcompresses these; full tools for backup pipes) ---
zstd
lz4
xz

# --- hardware identification ---
pciutils
usbutils
# cyme: modern lsusb replacement — readable USB tree with power/speed detail
cyme
dmidecode
lshw
lm-sensors
lm-sensors-detect

# --- network diagnostics ---
ethtool
iperf3
tcpdump
mtr
bind-tools

# --- monitoring ---
btop
# bottom (btm): btop-alternative system monitor with per-process I/O
bottom
iotop
ncdu
# duf: a better df — per-mount usage in a readable table (built locally from
# muesli/duf, repackaged under our signing key; NOT in Alpine v3.24)
duf
sysstat
fastfetch

# --- containers ---
docker
docker-cli-compose

# --- time / power / privilege ---
chrony
acpid
doas

# --- UPS monitoring (NUT); services off until configured ---
nut
nut-openrc
nut-udev

# --- outbound email for alerts (smartd, cron/SnapRAID output) ---
# msmtp is a send-only SMTP client (no daemon, no queue); mailx provides the
# mail(1) command smartd's warning script calls, wired to msmtp via the seeded
# /etc/mail.rc. Configure /etc/msmtprc first — README "Email alerts".
msmtp
mailx

# --- device firmware (curated consumer-x86 set, ~365 MB vs ~756 MB full) ---
# Target hardware = repurposed laptops / desktops / NUCs / mini-PCs from the
# last ~15 years: GPU, wifi, Bluetooth, wired-NIC, and laptop-platform
# firmware. Server NICs (mellanox/qed/cxgb*), ARM SoCs (qcom/mrvl/...), DSL,
# and TV capture are deliberately NOT shipped — every drop shrinks the
# modloop, the image, and boot time. Anything missing can be added on a
# running box (README: "Adding firmware for other hardware") — apk-added
# firmware installs early in boot, before devices are probed.
#
# GPU / display
linux-firmware-i915
linux-firmware-xe
linux-firmware-amdgpu
linux-firmware-radeon
linux-firmware-nvidia
# wifi (linux-firmware-intel = Intel wifi iwlwifi + Intel Bluetooth ibt)
linux-firmware-intel
linux-firmware-mediatek
linux-firmware-ath10k
linux-firmware-ath11k
linux-firmware-ath12k
linux-firmware-ath6k
linux-firmware-ath9k_htc
linux-firmware-brcm
linux-firmware-cypress
linux-firmware-rtw88
linux-firmware-rtw89
linux-firmware-rtlwifi
# Bluetooth (non-Intel)
linux-firmware-qca
linux-firmware-rtl_bt
linux-firmware-ar3k
# wired NICs found in desktops / mini-PCs
linux-firmware-rtl_nic
linux-firmware-tigon
linux-firmware-bnx2
linux-firmware-e100
# laptop platform (speaker amps, AMD platform/NPU, vendor quirks)
linux-firmware-cirrus
linux-firmware-amd
linux-firmware-amdnpu
linux-firmware-dell
linux-firmware-hp
linux-firmware-lenovo
linux-firmware-synaptics

# --- initramfs builder ---
mkinitfs

# --- OpenRC service scripts (modern Alpine ships init scripts in -openrc
#     subpackages; without these the /etc/init.d/* files don't exist and the
#     enabled services silently never start). The CI build verifies every
#     enabled service has an init script — adjust this list if the check fails.
busybox-openrc
openssh-server-common-openrc
dbus-openrc
avahi-openrc
chrony-openrc
acpid-openrc
rpcbind-openrc
smartmontools-openrc
docker-openrc
samba-server-openrc
nfs-utils-openrc
# optional services (enable per host; their script subpackage is here so it works):
lvm2-openrc
mdadm-openrc
tailscale-openrc
