# mountnas-tools/files/mountnas
#!/sbin/openrc-run
description="MountNAS storage guard and data-service supervisor"
DATA=/mnt/nasdata
STATE=/run/mountnas
BOOTMNT=/media/mnasboot
LOG=/var/log/mountnas.log
DATA_SERVICES="docker samba nfs"

depend() { need localmount; after net; }

log() { printf '%s mountnas: %s\n' "$(date '+%F %T')" "$*" >> "$LOG" 2>/dev/null; }
# fstab fstype for a mountpoint; "" if not declared
_fstype() { awk '$1!~/^#/ && $2=="'"$1"'"{print $3; exit}' /etc/fstab; }
_is_network_fs() { case "$1" in nfs|nfs4|cifs|smbfs|fuse.*|sshfs|glusterfs|ceph) return 0;; *) return 1;; esac; }

start() {
	mkdir -p "$STATE"; : > "$LOG" 2>/dev/null || true
	log "storage check starting"

	# Seed the prompt's unsaved-count cache. The PS1 indicator only READS this
	# file (never runs lbu itself), so the interactive prompt stays instant.
	lbu status 2>/dev/null | grep -c . > "$STATE/unsaved" 2>/dev/null || echo 0 > "$STATE/unsaved"

	# config partition (saves persist here). mountnas-mkdirs created /cfg and
	# localmount already tried this entry; we retry only as a safety net if the
	# device was not ready that early. /cfg is a LOCAL ext4 partition, so this
	# retry cannot hang the way a remote mount can.
	mountpoint -q /cfg || { mkdir -p /cfg; mount /cfg 2>>"$LOG" || log "/cfg mount failed"; }
	if mountpoint -q /cfg; then echo ok > "$STATE/cfg"; else echo fail > "$STATE/cfg"
		ewarn "Config partition /cfg NOT mounted — 'nas commit' will refuse to save."; log "/cfg NOT mounted"; fi

	# Bind the on-media local apk repo to a stable path for offline apk.
	# The '! mountpoint -q' guard makes this idempotent: on 'mountnas restart' a
	# second bind would STACK another mount on the same target instead of
	# replacing it, slowly leaking mounts. Bind only when nothing is there yet.
	# mkdir first: the mountpoint lives on the RAM root and exists on no boot —
	# without it the mount fails ("mount point does not exist"), the bind never
	# happens, and offline apk breaks (found via a live-box diagnostic).
	mkdir -p "$BOOTMNT"
	if mountpoint -q "$BOOTMNT" || mount -o ro LABEL=BOOT "$BOOTMNT" 2>>"$LOG"; then
		mkdir -p "$STATE/apks"
		if [ -d "$BOOTMNT/apks" ] && ! mountpoint -q "$STATE/apks"; then
			mount --bind "$BOOTMNT/apks" "$STATE/apks" 2>>"$LOG" || log "apks bind failed"
		fi
	fi

	# apk cache lives on the config partition (/cfg/cache, next to the apkovl;
	# /etc/apk/cache symlinks there). Ensure it exists, then sync the installed
	# set to /etc/apk/world: a user-added package is not in the on-media repo, so
	# the diskless init may have skipped it — it installs now from the cache (or
	# the network repos when up). This is what makes 'apk add X && nas commit'
	# stick across reboots. A no-op solver run when nothing is missing.
	if mountpoint -q /cfg; then
		mkdir -p /cfg/cache 2>>"$LOG" || true
		# shellcheck disable=SC2046  # world entries are single tokens by format
		apk add --quiet $(cat /etc/apk/world 2>/dev/null) >>"$LOG" 2>&1 \
			|| log "world sync (apk add) reported errors — a package may be missing until network/cache provides it"
	fi

	# Protect declared LOCAL data mounts that failed to mount, with a read-only
	# tmpfs placeholder (browsable but unwritable). Skip /mnt/nasdata (handled
	# specially below). Skip network filesystems entirely: we must never try to
	# mount them (a dead remote would hang boot) nor placeholder them (that would
	# clobber a mountpoint that may come up later or be written by a daemon).
	awk '$1!~/^#/ && $2 ~ /^\/mnt\//{print $2}' /etc/fstab | while read -r mp; do
		[ "$mp" = "$DATA" ] && continue
		ft=$(_fstype "$mp"); _is_network_fs "$ft" && continue
		if mountpoint -q "$mp"; then continue; fi
		mkdir -p "$mp"
		if mount -t tmpfs -o ro,size=4k,mode=0555 mountnas-blocked "$mp" 2>>"$LOG"; then
			log "placeholder (ro tmpfs) over unmounted $mp"
		fi
	done

	# data tier (system disk, always local)
	if ! awk '$1!~/^#/ && $2=="'"$DATA"'"{f=1} END{exit !f}' /etc/fstab; then
		echo fresh > "$STATE/data"; log "no data disk configured"
		ewarn "No data disk configured. Add it to /etc/fstab, then: nas commit"
		ewarn "Docker/Samba/NFS stay off until $DATA is mounted."
		return 0
	fi
	# A network filesystem as the system disk is unsupported: a dead remote
	# would stall this mount — and with it the whole default runlevel, which
	# busybox init finishes BEFORE it starts the inittab gettys, so a hung mount
	# here means no console login at all. Same guard the placeholder loop above
	# applies to the other /mnt entries.
	ft=$(_fstype "$DATA")
	if _is_network_fs "$ft"; then
		echo netfs > "$STATE/data"; log "$DATA declared as network fs ($ft) — unsupported, services held"
		eerror "$DATA is a network filesystem ($ft) — unsupported for the system disk; use a local disk. Services held."
		return 0
	fi
	# retry /mnt/nasdata once (local disk; safe, no hang). No placeholder exists
	# here — the loop above skipped it and stop() clears placeholders on restart.
	if ! mountpoint -q "$DATA"; then
		mkdir -p "$DATA"; mount "$DATA" 2>>"$LOG" || true
	fi
	if mountpoint -q "$DATA"; then
		mkdir -p "$DATA/docker" "$DATA/appdata" "$DATA/backups"
		echo ok > "$STATE/data"; log "$DATA mounted; starting data services"
		einfo "$DATA mounted — starting data services."
		for s in $DATA_SERVICES; do
			rc-service "$s" status >/dev/null 2>&1 || rc-service "$s" start 2>>"$LOG" || log "$s failed to start"
		done
		return 0
	fi
	dev=$(awk '$1!~/^#/ && $2=="'"$DATA"'"{print $1; exit}' /etc/fstab)
	if [ -n "$dev" ] && ! findfs "$dev" >/dev/null 2>&1; then
		echo disconnected > "$STATE/data"; log "$DATA device $dev NOT FOUND"
		eerror "Data disk $DATA ($dev) NOT FOUND — check power/cabling. Services held."
	else
		echo mountfail > "$STATE/data"; log "$DATA present but mount failed"
		eerror "Data disk $DATA present but FAILED to mount. Services held."
	fi
	# block stray writes to the empty system mountpoint, too
	mkdir -p "$DATA"; mount -t tmpfs -o ro,size=4k,mode=0555 mountnas-blocked "$DATA" 2>>"$LOG" || true
	return 0     # always succeed; state is reported, never a failed boot service
}

stop() {
	for s in $DATA_SERVICES; do
		rc-service "$s" status >/dev/null 2>&1 && rc-service "$s" stop
	done
	# drop any read-only placeholders so a restart can retry the real mounts
	awk '$1=="mountnas-blocked"{print $2}' /proc/mounts | while read -r mp; do
		umount "$mp" 2>/dev/null || true
	done
	return 0
}
