# mountnas-tools/files/mountnas-ttyd
#!/sbin/openrc-run
# mountnas-ttyd — browser-based terminal (ttyd), OFF by default.
# 'nas ttyd on' is the opt-in switch (same pattern as mountnas-web).
#
# Serves a REAL /bin/login prompt — never a bare shell: credentials are
# required, exactly like the physical console. Note that /etc/securetty has
# no pts/* entries, so ROOT cannot log in here — use a wheel user
# (adduser <name> wheel), which is the right posture for a browser terminal.
# Plain HTTP on the trusted LAN (credentials transit in cleartext on the
# wire — same trust model as Samba); 'nas ttyd on' says so out loud.
#
# Port defaults to 22222; override via PORT= in /etc/conf.d/mountnas-ttyd
# ('nas ttyd on <port>' writes it).
description="MountNAS browser terminal (ttyd login prompt)"

PORT=${PORT:-22222}
PIDFILE=/run/mountnas-ttyd.pid

depend() {
	need localmount
	after net
}

start_pre() {
	if ! command -v ttyd >/dev/null 2>&1; then
		eerror "ttyd not found — install it with: apk add ttyd"
		return 1
	fi
}

start() {
	ebegin "Starting MountNAS browser terminal (port $PORT)"
	# -W: writable (ttyd >=1.7 defaults read-only); the client gets a login
	# prompt on a fresh pty per connection
	start-stop-daemon --start --background \
		--make-pidfile --pidfile "$PIDFILE" \
		--exec /usr/bin/ttyd -- -W -p "$PORT" -t titleFixed="MountNAS terminal" /bin/login
	eend $?
}

stop() {
	ebegin "Stopping MountNAS browser terminal"
	start-stop-daemon --stop --quiet --pidfile "$PIDFILE" 2>/dev/null
	rm -f "$PIDFILE"
	eend 0
}
