# mountnas-tools/files/mountnas-net
#!/sbin/openrc-run
description="MountNAS first-boot network bootstrap"
depend() { need localmount; before networking; }

start() {
	# Do nothing if the user has configured any non-lo interface already.
	if awk '/^[[:space:]]*(auto|iface)[[:space:]]+/ && $2!="lo"{found=1} END{exit !found}' /etc/network/interfaces; then
		return 0
	fi
	# Pick the first real wired NIC (prefer one with carrier) — shared helper,
	# the same selection 'nas setup' uses.
	pick=$(/usr/libexec/mountnas/pick-nic 2>/dev/null) || pick=""
	[ -n "$pick" ] || { ewarn "no wired NIC found to auto-configure"; return 0; }
	einfo "First-boot: configuring DHCP on $pick (edit /etc/network/interfaces to customize)"
	cat >> /etc/network/interfaces <<EOF

auto $pick
iface $pick inet dhcp
EOF
	return 0
}
